# cats.service.template
[Unit]
Description=Cats website

[Service]
WorkingDirectory=/path/to/cats/repo
Type=simple
ExecStart=/path/to/cats/repo/.env/bin/python cats.py
User=gagarski

[Install]
WantedBy=multi-user.target